# systemd/aup-watchdog.service.example
# Example systemd user unit for the claude-tgbot AUP/idle/turns watchdog.
# Copy to the CONSUMER repo alongside the bot service, substitute values, install as:
#   ln -s $CONSUMER/systemd/claude-<consumer>-aup-watchdog.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now claude-<consumer>-aup-watchdog.service

[Unit]
Description=<CONSUMER_NAME> claude-tgbot AUP / idle / turns watchdog
Documentation=https://github.com/vjt/claude-tgbot
After=default.target

[Service]
Type=simple
# Consumer cwd used for: jsonl project-dir lookup, pane-cwd match, scrub_prompt.txt read.
Environment=CLAUDE_TGBOT_CONSUMER_DIR=<CONSUMER_DIR>
# Optional: Telegram chat_id to DM via <CONSUMER_DIR>/bot.send (SAY verb) when
# the watchdog can't resolve a claude pane for RESOLVE_ALERT_SEC (default 5 min).
# Unset = silent (journal-only).
#Environment=CLAUDE_TGBOT_ESCALATE_CHAT=<CHAT_ID>
WorkingDirectory=<CONSUMER_DIR>
ExecStart=/usr/bin/python3 -u /path/to/claude-tgbot/aup_watchdog.py
Restart=on-failure
RestartSec=5
# stdout/stderr default to journal; view with
#   journalctl --user -u claude-<consumer>-aup-watchdog.service -f

[Install]
WantedBy=default.target
