# systemd/consumer.service.example
# Example systemd user unit for a claude-tgbot consumer.
# Copy to the CONSUMER repo (not here), substitute paths, install as:
#   ln -s $CONSUMER/systemd/claude-tgbot-bot.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now claude-tgbot-bot.service

[Unit]
Description=<CONSUMER_NAME> Telegram bridge (claude-tgbot driving <CONSUMER_DIR>)
Documentation=https://github.com/vjt/claude-tgbot
After=network-online.target default.target
Wants=network-online.target

[Service]
Type=simple
Environment=CLAUDE_TGBOT_CONSUMER_DIR=<CONSUMER_DIR>
WorkingDirectory=<CONSUMER_DIR>
ExecStart=/path/to/.venv-tgbot/bin/python -u /path/to/claude-tgbot/bot.py
Restart=always
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=append:<CONSUMER_DIR>/bot.stdout.log
StandardError=inherit

[Install]
WantedBy=default.target
